# tasks/0.0.2/templates/mosaic-nebula-sync.timer.j2
[Unit]
Description=Sync Nebula config every 10 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
