# tasks/0.0.2/artifacts/mosaic-nebula-sync.service.j2
[Unit]
Description=Sync Nebula config

[Service]
Type=oneshot
ExecStart=mosaic-nebula sync
